# build-deb/litegix-agent_1.0-1_amd64/etc/systemd/system/litegix-agent.service
[Service]
ExecStart=/litegix/litegix-agent/litegix
Restart=always
StandardOutput=syslog
StandardError=syslog
User=root
Group=root
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
